# .config/systemd/user/network-office.service
[Install]
WantedBy=basic.target

[Unit]
Description=Office VPN Gateway

[Service]
ExecStart=/home/operator/bin/network-office
Restart=always
