# rootfs/etc/udev/rules.d/88-kbd.rules
ACTION=="add", SUBSYSTEM=="input", ENV{ID_INPUT_KEYBOARD}=="1", RUN+="/home/operator/bin/xkb-udev"
