# rootfs/etc/systemd/system/netconsoler.service
[Unit]
Description=netconsole receiver

[Service]
ExecStart=/bin/socat -u udp4-recv:6666 stdout
Restart=always

[Install]
WantedBy=basic.target
