# rootfs/etc/systemd/system/netconsoler.service
[Unit]
Description=netconsole receiver

[Service]
ExecStart=/bin/socat udp-recv:6666 -
Restart=always

[Install]
WantedBy=basic.target
