# root/etc/udev/rules.d/88-reconfig-kbd.rules
ACTION=="add", SUBSYSTEM=="usb", ENV{DISPLAY}=":0", ENV{XAUTHORITY}="/home/vxw/.Xauthority", RUN+="/home/vxw/bin/reconfig-kbd-udev"
